# deploy/supervisor/supervisord.d/example.ini
[program:Sms]

command=dotnet Sudo.Sms.WebApi.dll ; 运行程序的命令
directory=/vdb1/www/sms            ; 命令执行的目录

environment=ASPNETCORE_ENVIRONMENT=Production ; 进程环境变量

user=zsitd ; 进程执行的用户身份
stopsignal=INT

autostart=true
autorestart=true ; 程序意外退出 是否 自动重启
startsecs=3      ; 自动重启间隔

stderr_logfile=/var/log/supervisor/WebApiDemo.err.log ; 错误日志文件
stdout_logfile=/var/log/supervisor/WebApiDemo.out.log ; 输出日志文件

